# server/configs/dev/workhorse/supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run//supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run//supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf

[program:celery_task_workers]
directory=/cloud/src/jumodjango
command=/cloud/virtual/bin/python manage.py celeryd -Q now -c 2 -f /cloud/logs/celery_tasks.log -l DEBUG
autostart=true
autorestart=true
environment=PYTHONPATH=/cloud/src/jumodjango
environment=DJANGO_SETTINGS_MODULE=settings
environment=PATH=/cloud/virtual/bin/:$PATH


[program:celery_deferred_task_workers]
directory=/cloud/src/jumodjango
command=/cloud/virtual/bin/python manage.py celeryd -Q deferred -c 10 -f /cloud/logs/celery_deferred_tasks.log -l DEBUG
autostart=true
autorestart=true
environment=PYTHONPATH=/cloud/src/jumodjango
environment=DJANGO_SETTINGS_MODULE=settings
environment=PATH=/cloud/virtual/bin/:$PATH

#[program:celery_billing_workers]
#directory=/cloud/src/jumodjango
#command=/cloud/virtual/bin/python manage.py celeryd -Q billing -c 2 -f /cloud/logs/celery_billing_tasks.log -l DEBUG
#autostart=true
#autorestart=true
#environment=PYTHONPATH=/cloud/src/jumodjango
#environment=DJANGO_SETTINGS_MODULE=settings
#environment=PATH=/cloud/virtual/bin/:$PATH

[group:celery]
programs=celery_task_workers, celery_deferred_task_workers
priority=998

#[program:rabbitmq]
#command=/usr/sbin/rabbitmq-server
#autostart=true
#autorestart=true
#priority=999

#[group:billing]
#programs=celery_billing_workers
#priority=998
